# roles/orchestrator/templates/orchestrator.service.j2
[Unit]
Description={{ service_description }}
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_home }}
ExecStart={{ app_bin_path }} -config {{ app_config_path }}/config.json
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ service_name }}

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ app_home }}

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
